# debian/postinst
#!/bin/sh
set -e

inject_script_tag() {
    TEMPLATE="/usr/share/pve-manager/index.html.tpl"
    if [ -f "$TEMPLATE" ] && ! grep -qF "s3storage.js" "$TEMPLATE"; then
        SCRIPT_TAG='    <script type="text/javascript" src="/pve2/js/s3storage.js"></script>'
        sed -i "/pvemanagerlib\.js/a\\${SCRIPT_TAG}" "$TEMPLATE"
    fi
}

case "$1" in
    configure)
        inject_script_tag

        systemctl daemon-reload
        systemctl enable proxs3d.service
        # Don't auto-start — user needs to configure first
        echo ""
        echo "ProxS3 installed. Edit /etc/proxs3/proxs3d.json then run:"
        echo "  systemctl start proxs3d"
        echo ""
        echo "Add S3 storage via the web UI (Datacenter -> Storage -> Add -> S3)"
        echo "or via the command line:"
        echo "  pvesm add s3 my-s3-store --endpoint s3.example.com --bucket my-bucket \\"
        echo "    --access-key AKID... --secret-key SECRET... --content iso,vztmpl,snippets"
        ;;
    triggered)
        # pve-manager was upgraded and overwrote index.html.tpl — re-inject
        inject_script_tag
        ;;
esac

#DEBHELPER#
exit 0
